# resources/sysup/linuxtoys-update.service
[Unit]
Description=LinuxToys Update Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env EASY_CLI=1 UPD_SERVICE=1 linuxtoys
